# autostart/pcloud.desktop
[Desktop Entry]
Type=Application
Version=1.0
Name=pcloud
Comment=pcloudstartup script
Exec=/usr/local/bin/pcloud
StartupNotify=false
Terminal=false